# linux-protect/linux-protect.service
[Unit]
Description=Linux Protect Anti-Malware Daemon
Documentation=https://github.com/Cairn/linux-protect
After=local-fs.target systemd-tmpfiles-setup.service
ConditionPathExists=/etc/linux-protect/rules
ConditionPathIsDirectory=/etc/linux-protect/rules

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p /run/linux-protect
ExecStart=/usr/local/bin/linux-protect -config /etc/linux-protect/linux-protect.toml -v
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# ── Privilege ──────────────────────────────────────────────────────────
# fanotify FAN_OPEN_EXEC_PERM + FAN_ACCESS_PERM requires CAP_SYS_ADMIN.
# FAN_ACCESS_PERM enables shared library (LD_PRELOAD) interception.
# SIGKILL delivery requires appropriate signal permissions.
User=root
Group=root

# ── Capability bounding ─────────────────────────────────────────────────
AmbientCapabilities=
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_KILL CAP_DAC_OVERRIDE CAP_DAC_READ_SEARCH

# ── Filesystem isolation ──────────────────────────────────────────────
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/run/linux-protect /var/log
PrivateTmp=false
PrivateDevices=false

# ── Process isolation ──────────────────────────────────────────────────
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectProc=invisible

# ── Network ────────────────────────────────────────────────────────────
# Daemon is fully air-gapped. No inbound or outbound connections.
RestrictAddressFamilies=AF_UNIX AF_NETLINK
DenySockets=socketpair accept4 connect

# ── Resource limits ───────────────────────────────────────────────────
LimitNOFILE=65536
MemoryMax=256M
CPUQuota=25%

# ── Logging ────────────────────────────────────────────────────────────
StandardError=journal+console
StandardOutput=journal
SyslogIdentifier=linux-protect

[Install]
WantedBy=multi-user.target
